# roles/pangu-storm/files/ui.ini
[program:ui]
command=/opt/apache-storm-0.9.4/bin/storm ui
user=fudi
autostart=false
autorestart=true
startsecs=3
environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/root/bin
environment=PATH=/usr/lib/jvm/java-1.7.0-openjdk-1.7.0.85-2.6.1.2.el7_1.x86_64:%(ENV_PATH)s
environment=PATH=/usr/lib/jvm/jre-1.7.0-openjdk-1.7.0.85-2.6.1.2.el7_1.x86_64:%(ENV_PATH)s
environment=PATH=/opt/apache-storm-0.9.4/bin:%(ENV_PATH)s
